# tools/16ton/mtp_adb/android-tools-5.1.1r36+git20160322/debian/android-tools-adbd-reboot-watcher.conf
# Watchdog observing android reboot property.
# This is used by adbd to signal user asked for reboot
# Give android 1 second to complete its reboot before
# taking down also Ubuntu

start on android-container sys.powerctl=reboot*

task

script
    VAL=$(env | grep sys.powerctl=)
    sleep 1
    case ${VAL##*=} in
        reboot,)
            /sbin/reboot
        ;;
        reboot,recovery)
            /sbin/reboot -f recovery
        ;;
        reboot,bootloader)
            /sbin/reboot -f bootloader
        ;;
        reboot,poweroff)
            /sbin/poweroff
        ;;
        reboot,*)
            /sbin/reboot
        ;;
    esac
end script
